# crontab
* 11 * * * bash /usr/local/bin/do_updatedb.sh
